# crontab
* * * * * php /space-manager-v2/artisan schedule:run >> /dev/null 2>&1
